# stages/monitoring/files/configs/wtfd.service
[Unit]
Description=WTF just happened?

[Service]
ExecStart=/usr/local/bin/wtfd

[Install]
WantedBy=multi-user.target
