# stages/services/files/etc/systemd/system/davfs2-mounts/var-storage-smallwastebox.automount
[Unit]
Description=Hetzner Storagebox automount
After=network-online.target
Wants=network-online.target

[Automount]
Where=/var/storage/smallwastebox
TimeoutIdleSec=300

[Install]
WantedBy=remote-fs.target
