# stages/services/files/etc/systemd/system/containers.target
[Unit]
Description=A special target for container-based software
After=multi-user.target
