# stages/services/files/etc/systemd/system/containers.target
[Unit]
Description=A special target for container-based software

[Install]
WantedBy=multi-user.target
