# stages/services/files/etc/systemd/system/dave.service
[Unit]
Description=Dave, a WebDAV server

[Service]
Type=simple
User=www-data
ExecStart=/usr/local/bin/dave

[Install]
WantedBy=multi-user.target
