# stages/services/files/etc/systemd/system/drone-runner-amd64.service
[Unit]
Description=Drone CI Runner

[Service]
User=git
ExecStart=/usr/local/bin/drone-runner-amd64

[Install]
WantedBy=multi-user.target
