# stages/early/files/iacd.service
[Unit]
Description=Generalized IaC with batteries

[Service]
ExecStart=/usr/local/bin/iacd

[Install]
WantedBy=multi-user.target
