# stages/services/includes/media/files/etc/udev/rules.d/100-floppy-change.rules
ACTION=="change", ATTRS{idVendor}=="0644", ATTRS{idProduct}=="0000", ENV{DISK_MEDIA_CHANGE}=="1", RUN+="/usr/local/bin/media_mount $env{DEVNAME}"
