# stages/services/includes/system/files/etc/systemd/system/birdhouse-cleanup.service
[Unit]
Description=Cleanup before start

[Service]
Type=oneshot
ExecStart=/usr/local/bin/birdhouse-cleanup.sh

[Install]
WantedBy=multi-user.target
