# stages/services/includes/system/files/etc/systemd/system/disable-hdmi.service
[Unit]
Description=Disables HDMI output

[Service]
Type=oneshot
ExecStart=/usr/bin/tvservice -o

[Install]
WantedBy=multi-user.target
