# stages/services/includes/system/files/etc/systemd/system/sonar.service
[Unit]
Description=A dummy MP3 sonar

[Service]
Type=simple
ExecStart=/usr/local/bin/sonar

[Install]
WantedBy=multi-user.target
