# stages/services/includes/media/files/etc/systemd/system/davfs2-mounts/var-lib-chatserver-secondary.automount
[Unit]
Description=
After=network-online.target
Wants=network-online.target

[Automount]
Where=/var/lib/chatserver-secondary
TimeoutIdleSec=300

[Install]
WantedBy=remote-fs.target
