# stages/services/includes/media/files/etc/systemd/system/davfs2-mounts/var-storage-wastebox.automount
[Unit]
Description=Home NAS on RPi4 automount
After=wg-quick@wg0.service
Wants=network-online.target

[Automount]
Where=/var/storage/wastebox
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
